# rasp_pi_files/lora-uploader.service
[Unit]
Description=LoRa serial → AWS uploader
Wants=network-online.target
After=network-online.target porter-tunnel.service

[Service]
Type=simple
User=starfarm
Group=dialout
Environment=HOME=/home/starfarm
Environment=PATH=/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/starfarm/code/lora-farm-radio

# Optional, for consistency if env.sh references them
Environment=PORTER_TUNNEL_HOST=127.0.0.1
Environment=PORTER_TUNNEL_PORT=8122

# No ExecStartPre wait loop: the Python handles reachability and queues data.
ExecStart=/usr/bin/env bash -lc '/home/starfarm/code/lora-farm-radio/rasp_pi_files/run_logger.sh'

Restart=always
RestartSec=5s
# Give it time to stop cleanly on shutdown if needed
TimeoutStopSec=15s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
